# files/etc/crontabs/root
* * * * * nping --icmp --icmp-type 8 --data-string "flag-0xdeadbeef" --count 1 --interface eth0.111 10.1.11.100
